# camera/birdwatcher_camera.service
[Unit]
Description=Launches birdwatcher camera

[Service]
Type=oneshot
ExecStart=/usr/bin/birdwatcher_camera.py
