# camera/birdwatcher_camera.timer
[Unit]
Description=Run birdwatcher.service every 5 minutes

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
